# backend/geonature-service.conf
[program:geonature2]
command = APP_PATH/gunicorn_start.sh
autostart=true
autorestart=true
stdout_logfile = /var/log/supervisor/gunicorn_supervisor.log
redirect_stderr = true
