# install/assets/geonature-reload.service
[Unit]
Description=Reload GeoNature - ${GEONATURE_APP_NAME}

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reload ${GEONATURE_APP_NAME}