# install/assets/geonature-reload.service
[Unit]
Description=Reload GeoNature

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reload geonature
