# install/assets/geonature-reload.path
[Path]
PathModified=${BASE_DIR}/config/geonature_config.toml

[Install]
WantedBy=multi-user.target
