# install/assets/geonature-service.conf
[program:geonature2]
command = {{APP_PATH}}/gunicorn_start.sh
autostart=true
autorestart=true
stdout_logfile = {{ROOT_DIR}}/var/log/supervisor.log
redirect_stderr = true
user = {{USER}}
